# config/supervisord/hdfs.conf
[program:hdfs]
command=/usr/local/hadoop/sbin/start-dfs.sh
autostart=true
autorestart=false
redirect_stderr=true
