# config/supervisord/yarn.conf
[program:yarn]
command=/usr/local/hadoop/sbin/start-yarn.sh
autostart=true
autorestart=false
redirect_stderr=true
